# sandbox/guest/rootfs/etc/init.d/sync-agent
#!/sbin/openrc-run
description="Host-guest file sync agent (virtio-console)"

command="/usr/local/bin/sync-agent"
command_args="-root /workspace -dev /dev/hvc0"
supervisor="supervise-daemon"
respawn_delay=2
output_log="/var/log/sync-agent.log"
error_log="/var/log/sync-agent.log"

depend() {
	need workspace
}

start_pre() {
	# virtio-console device must exist; give the driver a moment on slow boots
	i=0
	while [ ! -e /dev/hvc0 ] && [ $i -lt 20 ]; do
		modprobe virtio_console 2>/dev/null
		sleep 0.5
		i=$((i + 1))
	done
	[ -e /dev/hvc0 ] || eerror "/dev/hvc0 missing"
	# raw mode: the framed binary protocol must not hit tty line discipline
	stty -F /dev/hvc0 raw -echo 2>/dev/null || true
}
