# sandbox/guest/etc/systemd/system/mount-share.service
# Format workspace disk if needed, mount /workspace, write WebDAV secrets.
# mount-share.sh manages /workspace lifecycle deterministically (formats only
# if blkid finds no filesystem, then mounts explicitly). workspace-sync.service
# depends on mount-share.service, not a systemd-fstab mount unit, eliminating
# a race where an unformatted first-boot disk would cause workspace.mount to
# fail permanently.
[Unit]
Description=Mount workspace disk and prepare WebDAV share
DefaultDependencies=no
After=dev-vdb.device
Requires=dev-vdb.device
Before=workspace-sync.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/libexec/mount-share.sh
ExecStop=/bin/sh -c 'umount /workspace 2>/dev/null || true'

[Install]
WantedBy=multi-user.target
